# isitdead-telegram.service
[Unit]
Description=Isitdead Telegram Integration Service
After=network.target isitdead.service
Wants=isitdead.service

[Service]
# Replace 'ubuntu' with your actual VPS username
User=ubuntu
Group=ubuntu

# Replace with the actual path to your project on the VPS
WorkingDirectory=/home/ubuntu/isitdead

# Path to the compiled Telegram integration binary
ExecStart=/home/ubuntu/isitdead/bin/isitdead-telegram

# Telegram integration listens on PORT, default is 8081.
# Required env vars:
# TOKEN=<telegram bot token>
# BASE_URL=<https://isitdead.cc>
# TELEGRAM_API_SECRET=<same value as main app>
# PORT=8081
EnvironmentFile=/home/ubuntu/isitdead/.env.telegram

# Automatic restart on crash
Restart=always
RestartSec=5

# Optional resource limits
MemoryMax=256M

[Install]
WantedBy=multi-user.target
